# kanardiacloud.service
[Unit]
Description=KanardiaCloud Flask Application
Documentation=https://github.com/rokmarko/kanardiacloud
After=network.target
Wants=network.target

[Service]
Type=exec
User=rok
Group=rok
WorkingDirectory=/home/rok/src/Cloud-1
Environment=PATH=/home/rok/src/Cloud-1/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=0
ExecStart=/home/rok/src/Cloud-1/start_production.sh
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kanardiacloud

# Security settings
NoNewPrivileges=true
ReadWritePaths=/home/rok/src/Cloud-1/instance
ReadWritePaths=/home/rok/src/Cloud-1/logs
ReadWritePaths=/var/log/kanardiacloud
ReadWritePaths=/var/run/kanardiacloud

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
